# systemd/rtcwake-daemon.service
[Unit]
Description=rtcwake planner daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/rtcwake-daemon
Restart=on-failure

[Install]
WantedBy=multi-user.target
